# pyinfra_borgbackup/borgbackup.service
[Unit]
Description=Run nightly backups
After=network.target

[Service]
Type=oneshot
ExecStart=/root/backup.sh


