# pyinfra_borgbackup/borgbackup.timer.j2
[Unit]
Description=Run nightly backups

[Timer]
OnCalendar=*-*-* {{ hour }}:{{ minute }}:00

[Install]
WantedBy=timers.target
